# root/etc/services.d/smartdns/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid root smartdns \
        -f \
        -c /config/smartdns.conf \
        -p /run/smartdns.pid
