# root/etc/services.d/smartdns/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc smartdns \
        -f \
        -c /config/smartdns.conf \
        -p /config/smartdns.pid
